# orchestrator/bootstrap/runasroot/rc.local
#!/bin/sh -e
sudo -Hiu airflow bash -c /opt/cassandra/bin/start_cassandra.sh
sudo -Hiu airflow bash -c /opt/hadoop/bin/start_hdfs.sh
sudo -Hiu airflow bash -c /opt/anaconda/bin/start_airflow.sh
exit 0
